# supervisord.conf
[supervisord]
nodaemon=true

[program:v2ray]
command=v2ray -config=/etc/v2ray/config.json
autorestart=true
priority=200
